# scripts/systemd/metrics_alerts.service
[Unit]
Description=Metrics Alerts Script (runs server_alerts_tracker.sh)
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/etc/default/metrics_bot_env
ExecStart=/bin/bash /path/to/repo/scripts/server_alerts_tracker.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
